# examples/calibration-lab/build/infra/calibration-lab-manager.container
# Generated by doql build — Podman Quadlet
[Unit]
Description=Calibration Lab Manager API
After=network-online.target

[Container]
Image=localhost/calibration-lab-manager:latest
PublishPort=8000:8000
EnvironmentFile=%h/.config/containers/calibration-lab-manager.env
Volume=calibration-lab-manager-data:/app/data:Z


[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
